# assesments/t1/backup-service/backup-tool.service
[Unit]
Description=Backup Tool Service
Documentation=man:backup-tool(1)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-tool.sh -c /etc/backup-tool/backup.conf
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
